# etc/profile.d/debian.sh
#!/bin/sh

#
# /etc/profile.d/debian.sh
#

# nala
which nala >/dev/null 2>&1 && alias apt='nala'

