# etc/profile.d/dotnet-cli-tools-bin-path.sh
#!/bin/sh

#
# /etc/profile.d/dotnet-cli-tools-bin-path.sh
#

append_path "${XDG_DATA_HOME}/dotnet/tools"

